# data_Crawl/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f

[program:api]
command=uvicorn dataCrawlAPI:app --host 0.0.0.0 --port 8001
